# vops.service.template
[Unit]
Description=vOps by vNodes[V] — infrastructure management and IP intelligence platform
After=network.target
Wants=network-online.target

[Service]
Type=simple
Environment=VOPS_HOME=__HOME__/.vOps
Environment=GOTRACEBACK=all
EnvironmentFile=-__HOME__/.vOps/.env
ExecStart=/usr/local/bin/vOps start
Restart=no
User=__VOPS_USER__
WorkingDirectory=__HOME__/.vOps
SyslogIdentifier=vOps
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
